# src/conf_supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:openvpn]
command=/bin/bash -c "cd /etc/openvpn; openvpn --config server_ovpn.conf"
